# conf/invoiceninja.cron
0 8 * * * __APP__ /usr/bin/php__YNH_PHP_VERSION__ __FINALPATH__/artisan ninja:send-invoices -q
0 8 * * * __APP__ /usr/bin/php__YNH_PHP_VERSION__ __FINALPATH__/artisan ninja:send-reminders -q

