# conf/invoiceninja.cron
0 8 * * * /var/www/invoiceninja/artisan ninja:send-invoices
0 8 * * * /var/www/invoiceninja/artisan ninja:send-reminders

